# crontab
0 6 * * * borg create -v --stats --compression lz4 /backup_repository::$(date +'%Y-%m-%d') /backup_sources
